# self-service/polysquid-self-service-nginx.service
[Unit]
Description=Self-Service HTTPS Reverse Proxy (Nginx)
After=docker.service network-online.target polysquid-self-service.service
Requires=docker.service
Wants=polysquid-self-service.service

[Service]
Type=simple
Restart=on-failure
RestartSec=10
TimeoutStopSec=30

ExecStartPre=-/usr/bin/docker rm -f polysquid_self_service_nginx

ExecStart=/usr/bin/docker run \
  --name polysquid_self_service_nginx \
  --add-host=host.docker.internal:host-gateway \
  -v /opt/polysquid/self-service/nginx/nginx.conf:/etc/nginx/nginx.conf:ro \
  -v /opt/polysquid/self-service/nginx/html/outside-subnet.html:/usr/share/nginx/html/outside-subnet.html:ro \
  -v /etc/polysquid/certs:/etc/nginx/certs:ro \
  -p 80:80 \
  -p 443:443 \
  nginx:1.27-alpine

ExecStop=/usr/bin/docker stop -t 20 polysquid_self_service_nginx
ExecStopPost=-/usr/bin/docker rm -f polysquid_self_service_nginx

[Install]
WantedBy=multi-user.target
